# service-mgmt/sm-common/debian/deb_folder/postinst
#!/bin/sh

set -e

mkdir -p /usr/lib/sm
mv /usr/lib64/libsm_common.so* /usr/lib/sm
systemctl enable sm-watchdog

#DEBHELPER#
exit 0