# service-mgmt/sm/debian/deb_folder/postinst
#!/bin/sh

set -e

/usr/bin/systemctl enable sm-shutdown.service >/dev/null 2>&1

#DEBHELPER#
exit 0